# test/docker/supervisor.conf
[supervisord]

nodaemon=true

[program:postgresql]

command = /etc/init.d/postgresql start
autorestart = true
