# test/docker/supervisor.conf
[supervisord]

nodaemon=true

[program:postgresql]

user=postgres
command = /usr/lib/postgresql/9.3/bin/postmaster -D /etc/postgresql/9.3/main
stopsignal=INT
autostart = true
autorestart = true
redirect_stderr=true
